# package/luci-app-aria2-bot-IRF/root/etc/uci-defaults/luci-aria2bot
#!/bin/sh
uci -q batch << EOB
  set aria2bot.main=aria2bot
  commit aria2bot
EOB

# Clear LuCI cache
rm -rf /tmp/luci-*
exit 0
