# package/luci-app-aria2-bot-IRF/files/etc/init.d/telegram_bot
#!/bin/sh /etc/rc.common
# Telegram Download Bot init script for OpenWrt
# luci-app-aria2-bot-IRF by irfanFRizki

START=99
STOP=10

USE_PROCD=1
PROG=/usr/bin/python3
BOT_SCRIPT=/root/telegram_download_bot.py
BOT_NAME="telegram_download_bot"
CONFIG_FILE=/etc/config/aria2bot

_info() {
	logger -p daemon.info -t "$BOT_NAME" "$*"
}

_err() {
	logger -p daemon.err -t "$BOT_NAME" "$*"
}

check_aria2_status() {
	if pgrep -x aria2c > /dev/null 2>&1; then
		return 0
	else
		return 1
	fi
}

check_aria2_rpc() {
	if command -v curl >/dev/null 2>&1; then
		curl -s -X POST http://localhost:6800/jsonrpc \
			-H "Content-Type: application/json" \
			-d '{"jsonrpc":"2.0","method":"aria2.getVersion","id":"1","params":[]}' \
			--max-time 3 >/dev/null 2>&1
		return $?
	fi
	return 0
}

wait_for_aria2() {
	local timeout=30
	local count=0
	_info "Waiting for Aria2 to be ready..."
	while [ $count -lt $timeout ]; do
		if check_aria2_status && check_aria2_rpc; then
			_info "Aria2 is ready"
			return 0
		fi
		sleep 1
		count=$((count + 1))
	done
	_err "Timeout waiting for Aria2"
	return 1
}

apply_config() {
	# Baca konfigurasi dari UCI dan patch ke bot script
	local bot_token allowed_users aria2_secret aria2_url max_concurrent max_retry

	config_load aria2bot
	config_get bot_token main bot_token ""
	config_get aria2_url main aria2_url "http://localhost:6800/jsonrpc"
	config_get aria2_secret main aria2_secret ""
	config_get max_concurrent main max_concurrent "2"
	config_get max_retry main max_retry "3"

	# Ambil allowed_users sebagai string
	allowed_users=""
	config_list_foreach main allowed_users handle_user() {
		local u="$1"
		if [ -z "$allowed_users" ]; then
			allowed_users="$u"
		else
			allowed_users="$allowed_users,$u"
		fi
	}

	if [ -z "$bot_token" ]; then
		_err "bot_token belum dikonfigurasi! Silakan set di LuCI > Services > Aria2 Bot > Settings"
		return 1
	fi

	# Patch bot script dengan konfigurasi UCI
	cp "$BOT_SCRIPT" "${BOT_SCRIPT}.run"

	# Replace token
	sed -i "s|TELEGRAM_BOT_TOKEN = \".*\"|TELEGRAM_BOT_TOKEN = \"${bot_token}\"|g" "${BOT_SCRIPT}.run"

	# Replace RPC URL
	sed -i "s|ARIA2_RPC_URL = \".*\"|ARIA2_RPC_URL = \"${aria2_url}\"|g" "${BOT_SCRIPT}.run"

	# Replace RPC Secret
	sed -i "s|ARIA2_RPC_SECRET = \".*\"|ARIA2_RPC_SECRET = \"${aria2_secret}\"|g" "${BOT_SCRIPT}.run"

	# Replace Max Concurrent
	sed -i "s|MAX_CONCURRENT_DOWNLOADS = [0-9]*|MAX_CONCURRENT_DOWNLOADS = ${max_concurrent}|g" "${BOT_SCRIPT}.run"

	# Replace Max Retry
	sed -i "s|MAX_RETRY_ATTEMPTS = [0-9]*|MAX_RETRY_ATTEMPTS = ${max_retry}|g" "${BOT_SCRIPT}.run"

	_info "Config applied: token=***${bot_token: -4}, rpc=${aria2_url}, concurrent=${max_concurrent}, retry=${max_retry}"
	return 0
}

start_service() {
	if [ ! -f "$BOT_SCRIPT" ]; then
		_err "Bot script tidak ditemukan di $BOT_SCRIPT"
		return 1
	fi

	if ! command -v python3 >/dev/null 2>&1; then
		_err "Python3 tidak ditemukan. Install: opkg install python3"
		return 1
	fi

	# Apply UCI config ke runtime script
	. /lib/functions.sh
	if ! apply_config; then
		return 1
	fi

	# Cek dan start aria2
	if ! check_aria2_status; then
		_info "Aria2 belum berjalan, mencoba start..."
		if [ -f /etc/init.d/aria2 ]; then
			/etc/init.d/aria2 start
			sleep 2
		else
			_err "Aria2 init.d tidak ditemukan. Install aria2 terlebih dahulu."
			return 1
		fi
	fi

	if ! wait_for_aria2; then
		_err "Aria2 RPC tidak merespons. Bot mungkin tidak berfungsi dengan baik."
	fi

	_info "Starting Telegram Download Bot..."

	procd_open_instance
	procd_set_param command "$PROG" "${BOT_SCRIPT}.run"
	procd_set_param respawn 3600 10 5
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param pidfile /var/run/${BOT_NAME}.pid
	procd_set_param file "$CONFIG_FILE"
	procd_close_instance

	_info "Telegram Download Bot started successfully"
}

stop_service() {
	_info "Stopping Telegram Download Bot..."
	rm -f "${BOT_SCRIPT}.run"
}

reload_service() {
	_info "Reloading Telegram Download Bot..."
	stop
	sleep 2
	start
}

service_triggers() {
	procd_add_reload_trigger "aria2bot"
}
