# package/luci-app-aria2-bot-IRF/files/etc/uci-defaults/99-luci-aria2bot
#!/bin/sh
# Post-install: enable bot service
/etc/init.d/telegram_bot enable
exit 0
